# .docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:php-fpm]
directory=/usr/share/nginx/html
command=php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

[program:cron]
command=cron -f
autostart=true
autorestart=true
stderr_logfile=/var/log/cron.err.log
stdout_logfile=/var/log/cron.out.log
user=root

[program:log-permission-fixer]
command=/bin/bash -c "while true; do find /usr/share/nginx/html/storage/logs -name '*.log' -not -user www-data -exec chown www-data:www-data {} \; 2>/dev/null || true; sleep 30; done"
autostart=true
autorestart=true
stderr_logfile=/var/log/log-permission-fixer.err.log
stdout_logfile=/var/log/log-permission-fixer.out.log
user=root

[program:log-monitor]
command=/bin/bash -c "while true; do /usr/share/nginx/html/.docker/monitor-logs.sh >> /var/log/log-monitor.log 2>&1; sleep 300; done"
autostart=true
autorestart=true
stderr_logfile=/var/log/log-monitor.err.log
stdout_logfile=/var/log/log-monitor.out.log
user=root

[program:artisan-migrate]
directory=/usr/share/nginx/html
command=/usr/bin/env bash -c "php artisan migrate --force"
autostart=true
autorestart=false
exitcodes=0
stdout_logfile=/var/log/artisan-migrate.log
stderr_logfile=/var/log/artisan-migrate.err.log
user=www-data
